# module06/configs/deep-health-check.service
[Unit]
Description=Deep Health Check for InfraScale
After=network.target
[Service]
Type=oneshot
ExecStart=/usr/local/bin/deep-health-check.sh
